# sunshine.service.in
[Unit]
Description=Sunshine Gamestream Server for Moonlight

[Service]
ExecStartPre=/bin/sh -c "test -d %E/sunshine || cp -r '@SUNSHINE_ASSETS_DIR@' '%E/sunshine'"
ExecStart=@SUNSHINE_EXECUTABLE_PATH@ %E/sunshine/sunshine.conf

[Install]
WantedBy=graphical-session.target
