# sunshine.service.in
[Unit]
Description=Sunshine Gamestream Server for Moonlight

[Service]
ExecStart=@SUNSHINE_EXECUTABLE_PATH@

[Install]
WantedBy=graphical-session.target
