# packaging/linux/sunshine.service.in
[Unit]
Description=@PROJECT_DESCRIPTION@
StartLimitIntervalSec=500
StartLimitBurst=5
PartOf=graphical-session.target
Wants=xdg-desktop-autostart.target
After=xdg-desktop-autostart.target

[Service]
ExecStart=@SUNSHINE_EXECUTABLE_PATH@
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=xdg-desktop-autostart.target
